# deployment/supervisor.ini
[program:flask_proj]
directory=/home/chris/webapps/flask_proj
command=/home/chris/webapps/flask_proj/.venv/bin/gunicorn -w 4 --bind unix:/home/chris/webapps/flask_proj/ftp.sock run:app
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data

